# docker/allinone/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=2
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; PostgreSQL 18 + pgvector
[program:postgres]
command=/usr/lib/postgresql/18/bin/postgres -D %(ENV_PGDATA)s
user=postgres
autostart=true
autorestart=true
priority=1
startsecs=3
startretries=5
stopsignal=SIGTERM
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres_err.log
environment=PGDATA="%(ENV_PGDATA)s"

; MinIO
[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address ":9001"
user=root
autostart=true
autorestart=true
priority=2
startsecs=3
startretries=5
stopsignal=SIGTERM
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/minio.log
stderr_logfile=/var/log/supervisor/minio_err.log

; Go 后端（等待 PG 就绪后启动）
[program:opsmind-server]
command=/bin/bash -c "while ! /usr/lib/postgresql/18/bin/pg_isready -h localhost -q; do sleep 1; done; exec /app/opsmind-server"
user=root
autostart=true
autorestart=unexpected
priority=3
startsecs=5
startretries=5
stopsignal=SIGTERM
stopwaitsecs=30
; logs go to supervisord stdout (visible in Railway)

; Next.js 前端（等待后端就绪后启动）
[program:opsmind-web]
command=/bin/bash -c "while ! curl -sf http://localhost:8080/health >/dev/null 2>&1; do sleep 1; done; exec node /app/web/server.js"
user=root
autostart=true
autorestart=unexpected
priority=4
startsecs=3
startretries=5
stopsignal=SIGTERM
stopwaitsecs=10
